# rpm/vpn-server-api.cron
RANDOM_DELAY=10

@hourly apache /usr/bin/vpn-server-api-housekeeping
@hourly apache /usr/bin/vpn-server-api-disconnect-expired-certificates
@daily apache /usr/bin/vpn-server-api-stats
