# rpm-dev/vpn-server-api.cron
RANDOM_DELAY=10

@hourly apache /usr/bin/vpn-server-api-housekeeping
@daily apache /usr/bin/vpn-server-api-stats
