# rpm-dev/vpn-user-portal.cron
RANDOM_DELAY=10

@hourly apache /usr/bin/vpn-user-portal-foreign-key-list-fetcher
